# other/pkgfile.cron
#!/bin/bash

. /etc/pkgtools/pkgfile.conf

if [[ $UPDATE_CRON ]]; then
    /usr/bin/pkgfile --update >/dev/null
fi
